# usr/lib/systemd/system/pve-fake-subscription.service
[Unit]
Description=Fake a Proxmox VE subscription

[Service]
Type=oneshot
ExecStart=/usr/bin/pve-fake-subscription
